# packages/botd/botd-14.tar.gz/botd-14/botd.service
[Unit]
Description=botd

[Service]
ExecStart=/usr/local/bin/botd
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target
